# ansible/roles/expand/files/rebootdelay.timer
[Unit]
Description=Delayed reboot

[Timer]
OnActiveSec=1
Unit=reboot.target
